# services/systemd/pulse-hub.service
[Unit]\nDescription=Pulse hub service\nAfter=network.target\n\n[Service]\nType=simple\nUser=pi\nWorkingDirectory=/opt/pulse/services/hub\nEnvironment=PYTHONUNBUFFERED=1\nEnvironment=CONFIG_FILE=/opt/pulse/config/config.yaml\nExecStart=/opt/pulse/.venv/bin/python /opt/pulse/services/hub/main.py\nRestart=always\nRestartSec=3\nStandardOutput=append:/var/log/pulse/hub.log\nStandardError=append:/var/log/pulse/hub.err\n\n[Install]\nWantedBy=multi-user.target\n