# services/systemd/pulse-health-monitor.service
[Unit]\nDescription=Pulse health monitor\nAfter=network.target\n\n[Service]\nType=simple\nUser=pi\nWorkingDirectory=/opt/pulse/services/sensors\nExecStart=/opt/pulse/.venv/bin/python /opt/pulse/services/sensors/health_monitor.py\nRestart=always\nRestartSec=3\nStandardOutput=append:/var/log/pulse/health_monitor.log\nStandardError=append:/var/log/pulse/health_monitor.err\n\n[Install]\nWantedBy=pulse-sensors.target\n