# services/systemd/pulse-hardware-detect.service
[Unit]\nDescription=Pulse hardware detect\nAfter=network.target\n\n[Service]\nType=oneshot\nUser=pi\nWorkingDirectory=/opt/pulse/services/sensors\nExecStart=/opt/pulse/.venv/bin/python /opt/pulse/services/sensors/hardware_detect.py\nStandardOutput=append:/var/log/pulse/hardware_detect.log\nStandardError=append:/var/log/pulse/hardware_detect.err\n\n[Install]\nWantedBy=multi-user.target\n