# services/systemd/pulse-dashboard.service
[Unit]
Description=Pulse dashboard web server
After=network.target pulse-hub.service
Wants=pulse-hub.service

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/pulse/dashboard/server
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=3
StandardOutput=append:/var/log/pulse/dashboard.log
StandardError=append:/var/log/pulse/dashboard.err

[Install]
WantedBy=multi-user.target
